# packaging/debian/prerm
#!/bin/sh
set -e
if [ "$1" = "remove" ] && command -v systemctl >/dev/null 2>&1; then
    systemctl disable --now caddy-l4.service || true
fi
